# etc/supervisor.d/valkey.conf
[program:valkey]
command=/usr/bin/valkey-server /etc/valkey/valkey.conf
user=root
startsecs=0
priority=1
autostart=true
autorestart=true
stopsignal=QUIT
